# rootfs/etc/s6-overlay/s6-rc.d/otbr-agent/run
#!/command/with-contenv bash
# ==============================================================================
# Standalone OpenThread Border Router Agent Run Script
# ==============================================================================
# This script replaces the original Home Assistant add-on script.
# It uses standard environment variables instead of bashio/supervisor API.
# ==============================================================================

set -e

# 1. Load Configuration from Environment
# ------------------------------------------------------------------------------
DEVICE="${DEVICE:-/dev/ttyUSB0}"
BAUDRATE="${BAUDRATE:-460800}"
FLOW_CONTROL="${FLOW_CONTROL:-1}"
OTBR_LOG_LEVEL="${OTBR_LOG_LEVEL:-info}"
FIREWALL="${FIREWALL:-1}"
NAT64="${NAT64:-1}"
BACKBONE_IF="${BACKBONE_IF:-}"

echo "[$(date)] INFO: Starting Standalone OTBR Agent..."
echo "[$(date)] INFO: Device: ${DEVICE}"
echo "[$(date)] INFO: Baudrate: ${BAUDRATE}"
echo "[$(date)] INFO: Flow Control: ${FLOW_CONTROL}"

# 2. Auto-detect Backbone Interface if not set
# ------------------------------------------------------------------------------
if [[ -z "${BACKBONE_IF}" ]]; then
    # Try to find the interface with the default route
    BACKBONE_IF=$(ip route | grep default | awk '{print $5}' | head -n1)
    if [[ -z "${BACKBONE_IF}" ]]; then
        echo "[$(date)] WARNING: Could not auto-detect backbone interface. Defaulting to eth0."
        BACKBONE_IF="eth0"
    else
        echo "[$(date)] INFO: Auto-detected backbone interface: ${BACKBONE_IF}"
    fi
else
    echo "[$(date)] INFO: Using configured backbone interface: ${BACKBONE_IF}"
fi

# 3. Verify Hardware
# ------------------------------------------------------------------------------
if [[ ! -e "${DEVICE}" ]]; then
    echo "[$(date)] ERROR: Device ${DEVICE} does not exist!"
    echo "[$(date)] ERROR: Available devices:"
    ls -la /dev/tty* 2>/dev/null || true
    # We exit here because OTBR cannot start without a radio
    exit 1
fi

if [[ ! -r "${DEVICE}" ]] || [[ ! -w "${DEVICE}" ]]; then
    echo "[$(date)] ERROR: Device ${DEVICE} is not readable/writable!"
    ls -la "${DEVICE}"
    exit 1
fi

# 4. Setup Firewall (re-implementation of upstream logic)
# ------------------------------------------------------------------------------
if [[ "${FIREWALL}" == "1" ]]; then
    echo "[$(date)] INFO: Configuring OTBR Firewall..."
    
    OTBR_FORWARD_INGRESS_CHAIN="OTBR_FORWARD_INGRESS"
    THREAD_IF="wpan0"
    
    # Create ipsets
    ipset create -exist otbr-ingress-deny-src hash:net family inet6
    ipset create -exist otbr-ingress-allow-dst hash:net family inet6
    
    # Flush existing rules
    ipset flush otbr-ingress-deny-src 2>/dev/null || true
    ipset flush otbr-ingress-allow-dst 2>/dev/null || true
    
    # Add multicast addresses to allow list
    ipset add -exist otbr-ingress-allow-dst ff02::1
    ipset add -exist otbr-ingress-allow-dst ff02::2
    ipset add -exist otbr-ingress-allow-dst ff03::1
    ipset add -exist otbr-ingress-allow-dst ff03::2
    ipset add -exist otbr-ingress-allow-dst ff03::fc
    
    # Configure ip6tables
    # Ensure chain exists
    if ! ip6tables -nL ${OTBR_FORWARD_INGRESS_CHAIN} >/dev/null 2>&1; then
        ip6tables -N ${OTBR_FORWARD_INGRESS_CHAIN}
    else
        ip6tables -F ${OTBR_FORWARD_INGRESS_CHAIN}
    fi
    
    # Insert jump rule if not present
    if ! ip6tables -C FORWARD -o ${THREAD_IF} -j ${OTBR_FORWARD_INGRESS_CHAIN} 2>/dev/null; then
        ip6tables -I FORWARD 1 -o ${THREAD_IF} -j ${OTBR_FORWARD_INGRESS_CHAIN}
    fi
    
    # Add rules to chain
    ip6tables -A ${OTBR_FORWARD_INGRESS_CHAIN} -m pkttype --pkt-type unicast -i ${THREAD_IF} -j DROP
    ip6tables -A ${OTBR_FORWARD_INGRESS_CHAIN} -m set --match-set otbr-ingress-deny-src src -j DROP
    ip6tables -A ${OTBR_FORWARD_INGRESS_CHAIN} -m set --match-set otbr-ingress-allow-dst dst -j ACCEPT
    ip6tables -A ${OTBR_FORWARD_INGRESS_CHAIN} -m pkttype --pkt-type unicast -j DROP
    ip6tables -A ${OTBR_FORWARD_INGRESS_CHAIN} -j ACCEPT
    
    echo "[$(date)] INFO: Firewall configured."
else
    echo "[$(date)] INFO: Firewall disabled by configuration."
fi

# 5. Start Dummy Syslog Listener (Crash Prevention)
# ------------------------------------------------------------------------------
# otbr-agent tries to connect to /dev/log. In Docker, this often doesn't exist.
# We create a dummy listener to keep it happy.
echo "[$(date)] INFO: Starting dummy syslog listener..."
python3 -c "
import socket, os, sys, signal, time

socket_path = '/dev/log'

def handler(signum, frame):
    sys.exit(0)

signal.signal(signal.SIGTERM, handler)

if os.path.exists(socket_path):
    try:
        os.remove(socket_path)
    except OSError:
        pass

try:
    s = socket.socket(socket.AF_UNIX, socket.SOCK_DGRAM)
    s.bind(socket_path)
    os.chmod(socket_path, 0o666)
    
    # Keep listener alive
    while True:
        try:
            data = s.recv(4096)
            # data is discarded (dummy listener)
        except Exception:
            pass
except Exception as e:
    print(f'Syslog setup failed: {e}')
" &
SYSLOG_PID=$!

# Wait for socket
timeout=50
while [ ! -S /dev/log ]; do
    sleep 0.1
    timeout=$((timeout - 1))
    if [ $timeout -le 0 ]; then
        echo "[$(date)] WARNING: Syslog socket failed to appear, otbr-agent might crash."
        break
    fi
done

echo "[$(date)] INFO: Syslog socket is ready."

# Small delay to ensure radio is ready after previous attempts
sleep 2

# 6. Construct Radio URL
# ------------------------------------------------------------------------------
RADIO_URL="spinel+hdlc+uart://${DEVICE}?uart-baudrate=${BAUDRATE}"

# Add flow control if enabled
if [[ "${FLOW_CONTROL}" == "1" ]]; then
    RADIO_URL="${RADIO_URL}&uart-flow-control"
else
    # Explicitly disable it if not enabled (though not adding the flag usually works)
    :
fi
# NOTE: We intentionally DO NOT add 'uart-init-deassert' (the fix)

echo "[$(date)] INFO: Radio URL: ${RADIO_URL}"

# 7. Start OTBR Agent
# ------------------------------------------------------------------------------
# We do NOT use -B (backbone) by default to prevent crashes on some hosts.
# We explicitly set -I wpan0 (default) and use -s for stdout logging.

echo "[$(date)] INFO: Executing otbr-agent..."
echo "[$(date)] INFO: Command: otbr-agent -I wpan0 -d${OTBR_LOG_LEVEL} -v -s ${RADIO_URL}"

# Run otbr-agent and capture both stdout and stderr
otbr-agent \
    -I wpan0 \
    -d"${OTBR_LOG_LEVEL}" \
    -v \
    -s \
    "${RADIO_URL}" 2>&1 &

OTBR_PID=$!
echo "[$(date)] INFO: otbr-agent started with PID ${OTBR_PID}"

# Wait for otbr-agent to finish
wait $OTBR_PID
EXIT_CODE=$?

if [ $EXIT_CODE -ne 0 ]; then
    echo "[$(date)] ERROR: otbr-agent exited with code ${EXIT_CODE}"
fi

exit $EXIT_CODE
