# rootfs/etc/s6-overlay/s6-rc.d/mdns/run
#!/bin/bash
# ==============================================================================
# Start mDNSResponder daemon
# ==============================================================================
echo "[$(date)] INFO: Starting mDNS Responder..."

# Start in debug mode (foreground) as expected by s6-rc
exec /usr/sbin/mdnsd -debug
